# runnerd/deploy/cambia-runnerd.service
[Unit]
Description=Cambia serving-harness runner daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cambia
WorkingDirectory=/srv/cambia
Environment=RUNNERD_BASE_DIR=/srv/cambia
Environment=RUNNERD_RUNS_DIR=/srv/cambia/runs
Environment=RUNNERD_MAX_CONCURRENT_JOBS=1
Environment=RUNNERD_MIN_FREE_RAM_GB=8
Environment=RUNNERD_MIN_FREE_DISK_GB=20
Environment=RUNNERD_JWT_PUBKEY=/srv/cambia/keys/jwt_ed25519.pub
Environment=RUNNERD_TLS_CERT=/srv/cambia/keys/tls.crt
Environment=RUNNERD_TLS_KEY=/srv/cambia/keys/tls.key
# v1 has no browser consumer; this gates browser-originated WS only. The CLI
# client sends no Origin header and is admitted on Bearer auth.
Environment=RUNNERD_ALLOWED_ORIGIN=https://dashboard.example.invalid
ExecStart=/usr/local/bin/cambia-runnerd --listen 0.0.0.0:8090
Restart=on-failure
RestartSec=5
NoNewPrivileges=yes
PrivateTmp=yes
KillMode=mixed
TimeoutStopSec=35

[Install]
WantedBy=multi-user.target
